# Lab3/task3.crontab
*/5 * * * 6 bash ~/home/user/Lab3/task1.sh
